# billing_cron_service/crontab
0 0 */3 * * . $HOME/.profile; python3 /opt/app/utils/put_users_to_queue.py >> /var/log/cron.log 2>&1
* * * * * . $HOME/.profile; python3 /opt/app/utils/yookassa_router.py >> /var/log/cron.log 2>&1
